# dist/init
#!/bin/sh
#
# ua2    This shell script takes care of starting and stopping the ua2 server
#
# chkconfig: 2345 80 30
# description: UA2 wastes your time
# processname: ICE 
# config: /etc/ua2/ua.edf
# pidfile: /var/run/ua2.pid

# Source function library.
. /etc/init.d/functions

# Source networking configuration.
[ -r /etc/sysconfig/network ] && . /etc/sysconfig/network

PID_FILE=/var/run/ua2.pid

start() {
    echo -n "Starting ua2: "
    daemonize -p ${PID_FILE} -u ua2 -c /var/lib/ua2 -o /var/log/ua2/ICE.stdout.log -- /usr/sbin/ICE -config /etc/ua2/ua.edf
    RETVAL=$?
    if [ $RETVAL -ne 0 ] ; then
        failure
    else
        success
    fi
    echo
}


stop() {
    echo -n "Stopping ua2: "
    killproc -p /var/run/ua2.pid ICE
    RETVAL=$?
    echo
}

restart() {
    stop
    start
}

rh_status() {
    status -p $PID_FILE ICE
}

rh_status_q() {
    rh_status >/dev/null 2>&1
}

# See how we were called.
case "$1" in
  start)
	rh_status_q && exit 0
        start
        ;;
  stop)
        stop
        ;;
  restart)
        restart
        ;;
  status)
        rh_status
        RETVAL=$?
        ;;
  *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 2
esac

exit $RETVAL
